# minv/package/profile.d/minv.sh
# minv initialization script (sh)
alias minv="sudo -u minv /usr/bin/minv"
